# dot_config/systemd/user/arch-maintenance.timer
[Unit]
Description=Run Arch Linux maintenance script weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
