# dot_config/systemd/user/arch-maintenance.service
[Unit]
Description=Run Arch Linux maintenance script

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/arch-maintenance.sh
